# mini-services/egxpy-bridge/egxpy-bridge.service
[Unit]
Description=EGXPy Bridge API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/egxpy-bridge
ExecStart=/opt/egxpy-bridge/venv/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8010 --workers 1 --log-level warning
Restart=always
RestartSec=5
Environment=PATH=/opt/egxpy-bridge/venv/bin:/usr/bin:/bin
Environment=EGXPY_PORT=8010
StandardOutput=append:/opt/egxpy-bridge/logs/service.log
StandardError=append:/opt/egxpy-bridge/logs/service-error.log

[Install]
WantedBy=multi-user.target
