# vps-service/egx-api.service
[Unit]
Description=EGX Data API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/egx-api
ExecStart=/usr/bin/python3 /opt/egx-api/egx_api_service.py
Restart=always
RestartSec=10
Environment=PORT=5000

[Install]
WantedBy=multi-user.target
